# apache/mautic.cron
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:leadlists:update >> /var/log/mautic.cron.log 2>&1
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:campaigns:update >> /var/log/mautic.cron.log 2>&1
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:campaigns:trigger >> /var/log/mautic.cron.log 2>&1
